# etc/supervisord.d/etherpad.ini
[program:etherpad]
command=/opt/qnib/etherpad/bin/start.sh
directory=/opt/etherpad/
user=etherpad
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
